# whoami/whoami.conf
[program:whoami]
user = "whoami"
group = "whoami"
command = "/usr/local/bin/whoami -port 8080"