# cloudflared/cloudflared.conf
[program:cloudflared]
user = cloudflared
group = cloudflared
command = /usr/local/bin/cloudflared tunnel run --token-file /etc/cloudflared/token