# cloudflared/cloudflared.conf
[program:cloudflared]
user = "cloudflared"
group = "cloudflared"
command = "/usr/local/bin/cloudflared tunnel run --token %(ENV_CF_TUNNEL_TOKEN)s"